# systemd/sudpd.service
[Unit]
Description=SUDP Server Daemon
After=network.target
Documentation=https://github.com/meirm/sudp

[Service]
Type=forking
ExecStart=/usr/local/bin/sudpd start --config-file /etc/sudp/server.yaml
ExecStop=/usr/local/bin/sudpd stop
ExecReload=/usr/local/bin/sudpd restart
PIDFile=/var/run/sudp/sudpd.pid

# Security
User=sudp
Group=sudp
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes

# Limits
LimitNOFILE=65535
LimitNPROC=4096

# Restart settings
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target 